# tools/profile.d/2_eip.sh
export PATH=$HOME/epics-train/tools/ether_ip/bin/$EPICS_HOST_ARCH:$PATH

