# tools/profile.d/3_example_iocs.sh

alias start_iocExample="cd ~/epics-train/examples/iocBoot/iocExample; ./st.cmd"

alias start_iocCombined="cd ~/epics-train/examples/iocBoot/iocCombined; ./st.cmd"


