# tools/profile.d/3_example_iocs.sh

alias start_iocExample="cd ~/epics-train/examples/iocBoot/iocExample; ./st.cmd"

alias start_iocCombined="cd ~/epics-train/examples/iocBoot/iocCombined; ./st.cmd"

alias start_imagedemo="~/epics-train/tools/EPICSV4Sandbox/ntndarrayServer/bin/$EPICS_HOST_ARCH/ntndarrayServerMain IMAGE"

alias start_neutrondemo="~/epics-train/tools/EPICSV4Sandbox/neutronsDemoServer/bin/$EPICS_HOST_ARCH/neutronServerMain -r"



