# tools/profile.d/4_css.sh

alias phoebus=~/epics-train/tools/start_phoebus.sh
alias css=~/epics-train/tools/start_phoebus.sh

alias start_iocfishtank="cd /home/training/epics-train/examples/Display\ Builder/fishtank; sh st.cmd"

