# tools/profile.d/4_python.sh
export PYTHONPATH=$HOME/epics-train/tools/pvaPy/lib/python/2.7/linux-x86_64

