# tools/profile.d/4_vdct.sh
alias vdct="/home/training/epics-train/tools/jdk-10/bin/java -cp ~/epics-train/tools/VisualDCT-2.8.0/VisualDCT.jar com.cosylab.vdct.VisualDCT ~/epics-train/tools/ether_ip/dbd/eipIoc.dbd"
