# tools/profile.d/4_vdct.sh
alias vdct=/home/training/epics-train/tools/run_vdct.sh
